# .config/etc/udev/rules.d/bluetooth-wakeup.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="8087", ATTRS{idProduct}=="0026" RUN+="/home/piotr/.local/script/kbrate.sh"
